# deploy/honestly-blog.service
[Unit]
Description=Honestly blog batch publish (3 passes per cycle, 24 articles total)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/honestly
EnvironmentFile=/opt/honestly/.env
ExecStart=/bin/bash /opt/honestly/deploy/honestly-batch-publish.sh
StandardOutput=append:/var/log/honestly-blog.log
StandardError=append:/var/log/honestly-blog.log
