# deploy/honestly-blog.timer
[Unit]
Description=Run the Honestly blog publish every 3 hours (8 refreshes/day, full district rotation)

[Timer]
# Every 3 hours starting at 01:07 UTC (off the :00 mark to avoid API stampedes).
# 8 cycles/day × 8 cities = 64 refreshes/day across 94 districts = full cycle in ~1.5 days.
OnCalendar=*-*-* 01:07:00
OnCalendar=*-*-* 04:07:00
OnCalendar=*-*-* 07:07:00
OnCalendar=*-*-* 10:07:00
OnCalendar=*-*-* 13:07:00
OnCalendar=*-*-* 16:07:00
OnCalendar=*-*-* 19:07:00
OnCalendar=*-*-* 22:07:00
Persistent=true

[Install]
WantedBy=timers.target
